# deploy/laravel-schedule.service
[Unit]
Description=Laravel scheduler (mailing-lists)
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/mailing-lists
ExecStart=/bin/php /var/www/mailing-lists/artisan schedule:work
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

#
#  ┌───────────────────────┬────────────────────────────────────────────────────────────────────────────────┐
#  │         Part          │                                     Reason                                     │
#  ├───────────────────────┼────────────────────────────────────────────────────────────────────────────────┤
#  │ schedule:work         │ Long-running process that calls schedule:run every minute — replaces system cron │
#  ├───────────────────────┼────────────────────────────────────────────────────────────────────────────────┤
#  │ Restart=always        │ Scheduler comes back if it crashes or the box reboots                          │
#  ├───────────────────────┼────────────────────────────────────────────────────────────────────────────────┤
#  │ After=mariadb.service │ Scheduled commands hit the DB — don't start before MariaDB is up               │
#  ├───────────────────────┼────────────────────────────────────────────────────────────────────────────────┤
#  │ User=frankenphp       │ Matches the owner of storage/ so scheduled commands can write logs/cache       │
#  └───────────────────────┴────────────────────────────────────────────────────────────────────────────────┘
#
# Install:
#   sudo cp deploy/laravel-schedule.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now laravel-schedule.service
#
# Inspect:
#   systemctl status laravel-schedule.service
#   journalctl -u laravel-schedule.service -f
